# infra/private-network/vps/linka-db-private-proxy.socket.template
[Unit]
Description=LINKa temporary private database cutover socket
After=wg-quick@wg0.service
Requires=wg-quick@wg0.service

[Socket]
ListenStream={{VPS_WG_ADDRESS}}:{{VPS_PRIVATE_PROXY_PORT}}
NoDelay=true

[Install]
WantedBy=sockets.target
