# infra/private-network/gateway/linka-private-firewall.service
[Unit]
Description=LINKa private gateway firewall
After=network-pre.target
Before=wg-quick@wg0.service
Wants=network-pre.target

[Service]
Type=oneshot
ExecStartPre=/bin/sh -c '/usr/sbin/nft list table inet linka_gateway >/dev/null 2>&1 || /usr/sbin/nft add table inet linka_gateway'
ExecStart=/bin/sh -c '{ /usr/bin/printf "flush table inet linka_gateway\n"; /usr/bin/cat /etc/linka-private-network/firewall.nft; } | /usr/sbin/nft --file -'
ExecReload=/bin/sh -c '{ /usr/bin/printf "flush table inet linka_gateway\n"; /usr/bin/cat /etc/linka-private-network/firewall.nft; } | /usr/sbin/nft --file -'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
